# playbooks/files/rescue-bot/rescue-bot.service
[Unit]
Description=KBTech out-of-cluster rescue bot
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service
ConditionPathExists=/etc/rescue-bot/rescue-bot.env

[Service]
Type=simple
User=kurt
WorkingDirectory=/opt/rescue-bot
EnvironmentFile=/etc/rescue-bot/rescue-bot.env
ExecStart=/opt/rescue-bot/venv/bin/python /opt/rescue-bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
